# pibOS/modules/pibbackend/filesystem/root/etc/systemd/system/doitpi_firstboot.service
[Unit]
Description=One time boot script
[Service]
Type=simple
ExecStart=/doitpi_firstboot.sh
[Install]
WantedBy=multi-user.target
